# crontab
* * * * * root python3 /tmp/update_config_files.py >> /var/log/cron.log